# .oblivious/services/server/run
#!/bin/sh
cd /code
python uvicorn_serve.py
